# dict_tests/files/usbmuxd.service
[Unit]
Description=Socket daemon for the usbmux protocol used by Apple devices
Documentation=man:usbmuxd(8)

[Service]
ExecStart=/usr/bin/usbmuxd --user usbmux --systemd
PIDFile=/var/run/usbmuxd.pid
